# src/webmount/webmountserver.service
[Unit]
Description=Isoft Web mount server

After=network.target

[Service]
ExecStart=/usr/sbin/webmountserver

[Install]
WantedBy=multi-user.target
